# pkgs/init/files/etc/inittab
::sysinit:/bin/init-new
::respawn:/etc/init.d/kubelet
::shutdown:/sbin/umount -a -r
