# pkgs/init/files/etc/inittab
::sysinit:/etc/init.d/rcS
::respawn:/etc/init.d/ip init
::respawn:/bin/busybox ntpd -n
::respawn:/bin/crio
::wait:/etc/init.d/init-node
::respawn:/etc/init.d/kubelet
::once:/etc/init.d/complete-init
