# pkgs/init/files/etc/inittab
::sysinit:/etc/init.d/rcS
::respawn:/etc/init.d/ip init
::respawn:/bin/busybox ntpd -np 0.pool.ntp.org
