# pkgs/init/files/etc/inittab
::sysinit:/etc/init.d/rcS
::respawn:/etc/init.d/kubelet
::shutdown:/sbin/umount -a -r
