# pkgs/init/files/etc/inittab
::sysinit:/etc/init.d/rcS
::respawn:/bin/crio
::respawn:/etc/init.d/kubelet
